# docker/supervisord.memory-worker.conf
; Slim supervisord profile for AURORACODER_ROLE=memory-worker (Layer 2b Gap
; Engine worker — see gateway/memory/ops/dispatcher.py). Only the stateless
; backend agent process runs: no Xvfb/VNC desktop, no gateway (no
; conversation/settings persistence needed for a one-shot investigation),
; no frontend, no toolstore-mgmt UI. Keeps the worker fast to boot and
; small — it is spawned on demand and torn down immediately after.
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:agent]
command=python3 run_web.py
directory=/app
environment=THINKTOOL_DOCKER="1",THINKTOOL_VNC="0"
autorestart=true
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
