# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:xvfb]
command=Xvfb :99 -screen 0 1280x720x24 -ac
autorestart=true
priority=10
stdout_logfile=/dev/null
stderr_logfile=/dev/null

[program:fluxbox]
command=fluxbox
environment=DISPLAY=":99"
autorestart=true
priority=20
stdout_logfile=/dev/null
stderr_logfile=/dev/null

[program:x11vnc]
command=x11vnc -display :99 -forever -nopw -shared -rfbport 5900
autorestart=true
priority=30
stdout_logfile=/dev/null
stderr_logfile=/dev/null

[program:novnc]
command=websockify --web /usr/share/novnc 6080 localhost:5900
autorestart=true
priority=40
stdout_logfile=/dev/null
stderr_logfile=/dev/null

[program:agent]
command=/opt/venv/bin/python run_web.py
directory=/app
environment=DISPLAY=":99",THINKTOOL_DOCKER="1",THINKTOOL_VNC="1"
autorestart=true
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:gateway]
command=/opt/venv/bin/python -m uvicorn gateway.api:app --host 0.0.0.0 --port 8081 --log-level warning
directory=/app
environment=BACKEND_URL="http://localhost:8080"
autorestart=true
priority=55
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:frontend]
command=/opt/venv/bin/python -m uvicorn frontend.server:app --host 0.0.0.0 --port 3000 --log-level warning
directory=/app
autorestart=true
priority=56
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:toolstore-mgmt]
command=/opt/venv/bin/python -m toolstore.webui.server --host 0.0.0.0 --port 8765
directory=/app
environment=TOOLSTORE_HOME="/app/data/toolstore"
autorestart=true
priority=57
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
